# Khmer Riel Currency forcesting using arima and lstm with web flask app throw docker/lstm/supervisord.conf
[supervisord]
nodaemon=true
[program:app]
command=python app_lstm.py --host=0.0.0.0
[program:api]
command=python app.py --host=0.0.0.0